# systemd/v2x-server.service
[Unit]
Description=V2X Dashboard Web Server - serves index.html/data.json over HTTP and reads STM32 telemetry over UART
After=network.target v2x-hub.service v2x-dashboard-bridge.service
Requires=v2x-hub.service

[Service]
Type=simple
User=rpi
WorkingDirectory=/home/rpi/V2X-Collision-Avoidance/RPI
ExecStart=/usr/bin/python3 /home/rpi/V2X-Collision-Avoidance/RPI/server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=v2x-server

[Install]
WantedBy=multi-user.target
