# finel systemd/v2x-server.service
[Unit]
Description=V2X Dashboard Web Server (serves index.html + data.json)
After=network.target v2x-hub.service v2x-dashboard-bridge.service
Requires=v2x-hub.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/rpi/v2x
ExecStart=/usr/bin/python3 /home/rpi/v2x/server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=v2x-server

[Install]
WantedBy=multi-user.target
