# finel systemd/v2x-hub.service
[Unit]
Description=V2X IPC Hub (Pub/Sub Broker)
# يبدأ قبل كل حاجة تانية في الشبكة
After=network.target
# كل الـ services التانية بتعتمد عليه
Before=v2x-dashboard-bridge.service v2x-car-client.service v2x-v2p.service v2x-server.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/rpi/v2x
ExecStart=/usr/bin/python3 /home/rpi/v2x/hub_.py
# لو وقع، ابدأه تاني بعد 3 ثواني
Restart=always
RestartSec=3
# سجل الـ output في الـ journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=v2x-hub

[Install]
WantedBy=multi-user.target
