# finel systemd/v2x-car-client.service
[Unit]
Description=V2X Car Client (MQTT Cloud + IPC Hub)
After=network-online.target v2x-hub.service v2x-dashboard-bridge.service
Requires=v2x-hub.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/v2x
ExecStart=/usr/bin/python3 /home/rpi/v2x/Car_client.py
Restart=always
RestartSec=5
# ننتظر شوية قبل أول تشغيل عشان الـ hub يبقى جاهز
ExecStartPre=/bin/sleep 2
StandardOutput=journal
StandardError=journal
SyslogIdentifier=v2x-car-client

[Install]
WantedBy=multi-user.target
