# IPC/IPC_Systemed/v2x-v2v-receiver.service
[Unit]
Description=V2X V2V UART Receiver (STM32 to RPi)
Documentation=V2X Collision Avoidance System
After=network.target v2x-hub.service
Requires=v2x-hub.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/v2x
ExecStart=/usr/bin/python3 /home/pi/v2x/v2v_uart_receiver.py

ExecStartPre=/bin/sleep 3

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=v2x-v2v-receiver

SupplementaryGroups=dialout

Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
