# IPC/IPC_Systemed/v2x-hub.service
[Unit]
Description=V2X IPC Hub (Pub/Sub Unix Socket Broker)
Documentation=V2X Collision Avoidance System
After=network.target
# Hub يجب أن يبدأ قبل كل الـ nodes التانية

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/v2x
ExecStart=/usr/bin/python3 /home/pi/v2x/hub.py

# إعادة التشغيل تلقائياً لو وقع
Restart=always
RestartSec=3

# لوج الـ output في journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=v2x-hub

# متغيرات البيئة
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
