# IPC_Systemed/v2x-v2p.service
[Unit]
Description=V2X V2P Camera Node (YOLO Pedestrian & Vehicle Detection)
Documentation=V2X Collision Avoidance System
After=network.target v2x-hub.service v2x-traffic-light.service
Requires=v2x-hub.service
# يحتاج الـ hub، والأفضل ينتظر الـ traffic light يبدأ أول

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/v2x
ExecStart=/usr/bin/python3 /home/pi/v2x/V2P_TEST.py

# انتظر أطول شوية عشان الكاميرا والموديل يحتاجوا وقت للتحميل
ExecStartPre=/bin/sleep 4

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=v2x-v2p

# الكاميرا محتاجة الـ video group
# تأكد إن الـ pi user في group video:
#   sudo usermod -aG video pi
SupplementaryGroups=video

Environment=PYTHONUNBUFFERED=1
Environment=QT_QPA_PLATFORM=offscreen

[Install]
WantedBy=multi-user.target
